# package/rp-pppoe/files/pppoe-client.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

wanproto=`uci get network.wan.proto`
START=50

start() {
	if [ "$wanproto" = "pppoe" ]; then
		pppoe-start
	fi
}

stop() {
	pppoe-stop
}

