# target/linux/ramips/base-files/etc/uci-defaults/leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/ramips.sh

board=$(ramips_board_name)

case "$board" in
wifi-router)
	ucidef_set_led_timer "system" "system" "demo:green:system" "550" "350"
	;;
wifi-router_v2)
	ucidef_set_led_timer "system" "system" "demo:green:system" "550" "350"
	;;
wifi-router-poe)
	ucidef_set_led_timer "system" "system" "demo:green:system" "550" "350"
	;;
wifi-router-wallpoe)
	ucidef_set_led_timer "system" "system" "demo:green:system" "550" "350"
	;;
esac

ucidef_commit_leds

exit 0
