# package/wifidog_ext/files/wifidog_ext.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=66

start() {
	/usr/bin/getwanip.sh
	/bin/upgrade_daemon &
	/bin/mythstart &
}

stop() {
	echo "running"
	#
}

