# package/luci/src/themes/bootstrap_ODM/root/etc/uci-defaults/luci-theme-bootstrap
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_ODM=/luci-static/bootstrap_ODM
	set luci.main.mediaurlbase=/luci-static/bootstrap_ODM
	commit luci
EOF
exit 0
